# files/llama-server.service
[Unit]
Description=LLama.cpp AI Server for Primaria Galati
After=network.target

[Service]
Type=simple
User=cristi
WorkingDirectory=/home/cristi/llama
ExecStart=/home/cristi/llama/llama-server \
  -m /home/cristi/llama/DeepSeek-V2-Lite-Chat-Q5_K_M.gguf \
  -ngl 9 \
  --ctx-size 2048 \
  --batch-size 128 \
  --threads 4 \
  --port 8080 \
  --host 0.0.0.0
Restart=on-failure
RestartSec=5
Environment=LD_LIBRARY_PATH=/home/cristi/llama
Environment=GGML_CUDA_ENABLE_UNIFIED_MEMORY=1

[Install]
WantedBy=multi-user.target
